# home/.local/sv/pipewire/run
#!/bin/sh

if [ -f ${SESSION_ENV} ]; then
	. ${SESSION_ENV}
else
	return 1
fi
export XDG_RUNTIME_DIR
export SVDIR

exec pipewire >/dev/null
