# modules/network/nethandler.service
[Unit]
Description=Nethandler

[Service]
ExecStart=/usr/local/bin/nethandler

[Install]
WantedBy=default.target

