# root/etc/init.d/ivan-runsvdir
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/bin/runsvdir
USER_SVCDIR="/home/ivan/.local/runsvdir"

export USER="ivan"
export HOME="/home/ivan"
export SESSION_ENV="/tmp/session.env"

command_user="${USER}" 
command_args="${USER_SVCDIR}"

start_pre()
{
	checkpath -d ${USER_SVCDIR}
}


stop_pre()
{
	sv force-shutdown ${USER_SVCDIR}/*
}
