# home/.local/sv/wireplumber/run
#!/bin/bash

if [ -f ${SESSION_ENV} ]; then
	. ${SESSION_ENV}
else
	return 1
fi
export XDG_RUNTIME_DIR
export SVDIR

sv check pipewire >/dev/null || exit 1
exec wireplumber >/dev/null
